# recipes-bsp/load-module/files/load-module.service
[Unit]
Description=load-module

[Service]
ExecStart=/usr/bin/load-module.sh

[Install]
WantedBy=multi-user.target
